# n8n/supervised.conf
[supervisord]
user=root
nodaemon=true
loglevel=debug

[program:nginx]
command=bash /app/entrypoint-nginx.sh
autorestart=true

[program:n8n]
command=bash /app/entrypoint-n8n.sh
autorestart=true